# drivers/fingerprint/fprintd-override.service
[Service]
ExecStart=
ExecStart=/opt/fingerpp/bin/fingerpp --bus=system --filter-before-ssim=true --min-score=0.3
